# discord-audit-stream.service
[Unit]
Description=Discord Audit Stream - OpenClaw session monitor
Documentation=https://github.com/Sabrimjd/discord-audit-stream
After=network.target

[Service]
Type=simple
User=sab
Group=sab
WorkingDirectory=/home/sab/.openclaw/hooks/discord-audit-stream
ExecStart=/usr/bin/node daemon.ts
Restart=on-failure
RestartSec=10

# Environment variables (optional - can also use config.json)
# Environment="DISCORD_AUDIT_WEBHOOK_URL=https://discord.com/api/webhooks/YOUR_ID/YOUR_TOKEN"
# Environment="DISCORD_AUDIT_CHANNEL_ID=YOUR_CHANNEL_ID"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=discord-audit-stream

[Install]
WantedBy=multi-user.target
